# systemd-global/leclient-restart.path
[Path]
PathChanged=/home/letsencrypt/.config/leclient/crt_renewed.trigger

[Install]
WantedBy=multi-user.target
